# tools/systemd/speech-to-text.service
[Unit]
Description=Speech-to-Text Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 %h/main/projects/speech-to-text/stt_daemon.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
Environment=PYTHONPATH=%h/main/projects/speech-to-text

[Install]
WantedBy=graphical-session.target
